# steam4j-mockapi/src/deb/upstart/steam4j-mockapi.conf
# vim: set ft=upstart ts=4 et:
description "steam4j-mockapi"

start on runlevel [2345]
stop on runlevel [!2345]

limit nofile 64000 64000

kill timeout 300 # wait 300s between SIGTERM and SIGKILL.

pre-start script
    mkdir -p /var/lib/steam4j-mockapi/
    mkdir -p /var/log/steam4j-mockapi/
end script

script
    ENABLE_STEAM4JMOCKAPI="yes"

    if [ "x$ENABLE_STEAM4JMOCKAPI" = "xyes" ]; then
        exec start-stop-daemon --start --quiet --chuid root  \
            --exec /opt/solong/steam4j-mockapi/current/bin/start -- 
    fi
end script
